# chatgpt-widget.service.template
[Unit]
Description=ChatGPT Desktop Widget
After=network.target

[Service]
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=YOUR_PROJECT_PATH
ExecStart=YOUR_NPM_PATH start
Restart=on-failure
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target 